# etl/main_project/etl_process/crontab
*/1 * * * * python3 etl_process.py >> /logfile.txt
